# cron.d/evaluate_submissions
0 * * * * root cd /app && export PATH="/app/.venv/bin:$PATH" && python manage.py evaluate_submissions >> /storage/cron.log 2>&1
